# .config/systemd/user/proton-drive-mount.service
[Unit]
Description=Mount Proton Drive via rclone
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemd-run --user --scope /bin/bash -c "%h/.dotfiles/.proton-drive-rclone-mount.sh -vv 2>&1 | tee /tmp/proton-drive-mount-$(date +%%Y%%m%%d-%%H%%M%%S).log"
Environment=DISPLAY=:0
Environment=XAUTHORITY=%h/.Xauthority
ExecStartPost=/usr/bin/sleep 10

[Install]
WantedBy=graphical-session.target

